# build/debian/mixxx.postrm.in
#!/bin/sh

# debhelper will insert its stuff here, e.g. update-rc.d-calls

#DEBHELPER#

exit 0
